# debian/dog-agent.service
[Unit]
Description=Dognet Unified Agent (FireDog / CyberSheppard / SentinelCore)
Documentation=https://github.com/dognet/dog-agent
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=60
StartLimitBurst=5

[Service]
Type=simple
User=dog-agent
Group=dog-agent
ExecStart=/usr/bin/dog-agent --config /etc/dog-agent/agent.conf
Restart=on-failure
RestartSec=5s

# Sicurezza
NoNewPrivileges=false
ProtectHome=true
ProtectSystem=strict
ReadWritePaths=/var/log/dog-agent /var/run/dog-agent
PrivateTmp=true
CapabilityBoundingSet=CAP_NET_ADMIN CAP_NET_RAW

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=dog-agent

[Install]
WantedBy=multi-user.target
