# debian/dog-agent.service
[Unit]
Description=Dognet Unified Agent (FireDog / CyberSheppard / SentinelCore)
Documentation=https://github.com/dognet/dog-agent
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=60
StartLimitBurst=5

[Service]
Type=simple
# L'agent deve poter manipolare iptables (lock file /run/xtables.lock è di root)
# e leggere /proc/<pid>/... per i collector di sistema. Gira come root con
# ProtectSystem=strict per limitarne i side-effect sul filesystem.
User=root
Group=root
ExecStart=/usr/bin/dog-agent --config /etc/dog-agent/agent.conf
Restart=on-failure
RestartSec=5s

# Sicurezza (anche da root, riduciamo la superficie scrivibile)
NoNewPrivileges=false
ProtectHome=true
ProtectSystem=strict
ReadWritePaths=/var/log/dog-agent /var/run/dog-agent /run/xtables.lock
PrivateTmp=true

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=dog-agent

[Install]
WantedBy=multi-user.target
